# discord-grok-bot.service
[Unit]
Description=Discord Grok Bot
Documentation=https://github.com/clotho2/wolfe-discord-bot-enhanced
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/aicara/wolfe-discord-bot-enhanced

# Environment file (contains all secrets)
EnvironmentFile=/opt/aicara/wolfe-discord-bot-enhanced/.env

# Start command - use full path to npm
ExecStart=/opt/node22/bin/npm start

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening (relaxed for compatibility)
NoNewPrivileges=true
PrivateTmp=true
# ProtectSystem=strict  # Disabled - can cause NAMESPACE errors
# ProtectHome=read-only  # Disabled - can cause NAMESPACE errors
# ReadWritePaths=/opt/aicara/wolfe-discord-bot-enhanced/logs  # Disabled - logs will use journal

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-grok-bot

# Resource limits (optional - adjust as needed)
# MemoryMax=2G
# CPUQuota=200%

[Install]
WantedBy=multi-user.target
