# discord-grok-bot.service
[Unit]
Description=Discord Grok Bot
Documentation=https://github.com/yourusername/discord-ai-bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/discord-ai-bot

# Environment file (contains all secrets)
EnvironmentFile=/opt/discord-ai-bot/.env

# Start command - run node directly (compiled JS is in src/ folder)
# Make sure to run 'npm run build' before starting the service
ExecStart=/root/.nvm/versions/node/v18.20.8/bin/node src/server.js

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=5min
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening (relaxed for compatibility)
NoNewPrivileges=true
PrivateTmp=true
# ProtectSystem=strict  # Disabled - can cause NAMESPACE errors
# ProtectHome=read-only  # Disabled - can cause NAMESPACE errors
# ReadWritePaths=/opt/discord-ai-bot/logs  # Disabled - logs will use journal

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-grok-bot

# Resource limits (optional - adjust as needed)
# MemoryMax=2G
# CPUQuota=200%

[Install]
WantedBy=multi-user.target
